# deploy/systemd/veda-server.service
# systemd unit for veda-server on a host with native MySQL + Milvus.
# Install:
#   sudo useradd --system --no-create-home --home-dir /data/veda veda
#   sudo install -d -o root -g root /data/veda /data/veda/bin /data/veda/config
#   sudo install -d -o veda -g veda /data/veda/data
#   sudo install -m 0755 target/release/veda-server /data/veda/bin/
#   sudo install -m 0640 -o root -g veda config/server.toml /data/veda/config/config.toml
#   sudo cp deploy/systemd/veda-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now veda-server.service
#
# Logs:  journalctl -u veda-server -f
# Reload after config edits:  sudo systemctl restart veda-server

[Unit]
Description=Veda server
Documentation=https://github.com/jooekong/veda
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=veda
Group=veda
# Optional env overrides — secrets here win over config.toml. The leading
# `-` makes the file optional so the unit still starts if it's absent.
EnvironmentFile=-/data/veda/config/veda.env
WorkingDirectory=/data/veda/data
ExecStart=/data/veda/bin/veda-server /data/veda/config/config.toml
Restart=on-failure
RestartSec=5s
# 10s shutdown grace lets in-flight HTTP requests finish before SIGKILL;
# the server's tokio shutdown handler already wires graceful drain.
TimeoutStopSec=10s

# Hardening. We only write under /data/veda/data so the rest is read-only.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/data/veda/data
# A 4 GB RSS cap catches a leak before it OOMs the whole box. Veda's
# steady-state RSS is well under 1 GB; raise this if you embed huge files.
MemoryMax=4G

[Install]
WantedBy=multi-user.target
