# deploy/systemd/veda-server.service
# systemd unit for veda-server on a host with native MySQL + Milvus.
# Install:
#   sudo cp deploy/systemd/veda-server.service /etc/systemd/system/
#   sudo useradd --system --no-create-home --home-dir /var/lib/veda veda
#   sudo install -d -o veda -g veda /etc/veda /var/lib/veda
#   sudo install -m 0640 -o veda -g veda config/server.toml /etc/veda/config.toml
#   sudo install -m 0755 target/release/veda-server /usr/local/bin/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now veda-server.service
#
# Logs:  journalctl -u veda-server -f
# Reload after config edits:  sudo systemctl restart veda-server

[Unit]
Description=Veda server
Documentation=https://github.com/jooekong/veda
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=veda
Group=veda
EnvironmentFile=-/etc/veda/veda.env
ExecStart=/usr/local/bin/veda-server /etc/veda/config.toml
Restart=on-failure
RestartSec=5s
# 10s shutdown grace lets in-flight HTTP requests finish before SIGKILL;
# the server's tokio shutdown handler already wires graceful drain.
TimeoutStopSec=10s

# Hardening. We don't write outside /var/lib/veda so the rest is read-only.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/veda
# A 4 GB RSS cap catches a leak before it OOMs the whole box. Veda's
# steady-state RSS is well under 1 GB; raise this if you embed huge files.
MemoryMax=4G

[Install]
WantedBy=multi-user.target
